# supervisord_gpu_optimized.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# 優化的 GPU Whisper 後端服務
[program:gpu-whisper-backend]
command=/usr/local/bin/python /app/gpu_whisper_server.py
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/gpu-whisper-backend.log
stderr_logfile=/var/log/supervisor/gpu-whisper-backend.log
environment=BACKEND_PORT="8080",CUDA_VISIBLE_DEVICES="0"
startsecs=30
startretries=3
stopwaitsecs=10

# nginx 前端代理
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
startsecs=5